# templates/default/dynamic-dynamodb.erb
#!/bin/bash
#
# chkconfig: 35 90 12
# description: Dynamic-Dynamodb init script
#
# Get function from functions library
. /lib/lsb/init-functions

PIDFILE=<%= @pidfile %>
# Start the service
start() {

        if [ -e $PIDFILE ]; then
                pgrep -f "<%= @command %>"
                if [ $? -ne 0 ]; then
                        echo "Cleaning up stale PID file: $PIDFILE"
                        rm $PIDFILE
                fi
        fi
        echo -n "echo -n Starting Dynamic-DynamoDB: "
        <%= @command %> --daemon start &
        echo $"Dynamic-DynamoDB Started"
}
# Restart the service dynamic-dynamodb
stop() {
        echo -n "echo -n Stopping Dynamic-DynamoDB: "
        <%= @command %> --daemon stop &
        echo $"Dynamic-DynamoDB Stopped"
}
### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        pgrep -f "<%= @command %>"
        ;;
  restart|reload|condrestart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac
exit $?
